# level_analyzer/deploy/level-analyzer.service
[Unit]
Description=Level Analyzer (Fase A) - notifiche zone conf=2 (Telegram)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=matteo
WorkingDirectory=/home/matteo/trad
ExecStart=/home/matteo/trad/venv/bin/python -m level_analyzer run
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
